# support/runit/log/run
#!/bin/sh
exec chpst -unobody svlogd /var/log/smartplug-monitor
